# monitor/statsd/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor 

[program:carbon]
command=nodejs /statsd/stats.js /statsd/config.js
autorestart=true
