# media/coturn/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor 

[program:coturn]
command=/usr/bin/turnserver -c /etc/turnserver.conf -v
autorestart=true
