# huao/mysqlrouter/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor 

[program:mysqlrouter]
command=/usr/sbin/mysqlrouter
autorestart=true
