# cloud/glusterfs/files/supervisord.conf
[supervisord]
nodaemon=true

[program:glusterd]
command=/usr/sbin/glusterd -N -p /usr/var/run/glusterd.pid -f /usr/etc/glusterfs/glusterd.vol
autorestart=true

