# monitor/graphite/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor 

[program:carbon]
command=/opt/graphite/bin/carbon-cache.py start --nodaemon
autorestart=true
